# supervisord.conf
[supervisord]
nodaemon=true
[program:nginx]
command=/user/bin/nginx -g daemon off
[program:tomcat]
command=/usr/share/tomcat7/bin/catalina.sh run
