# supervisord.conf
[supervisord]
nodaemon=true
[program:nginx]
command=/usr/sbin/nginx & 
[program:tomcat]
command=/usr/share/tomcat7/bin/catalina.sh run
[program:sh]
command=/bin/bash
